# unbound/rootfs/etc/services.d/unbound/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Unbound..."

CONFIG_DIR="/config/unbound"

# Create config directory in /config (Home Assistant config folder)
if [ -d "/config" ] && [ -w "/config" ]; then
    mkdir -p "${CONFIG_DIR}"

    # Create example config if no config exists yet
    if [ ! -f "${CONFIG_DIR}/unbound.conf" ] && [ ! -f "${CONFIG_DIR}/unbound.conf.example" ]; then
        if [ -f "/etc/unbound/unbound.conf.example" ]; then
            cp /etc/unbound/unbound.conf.example "${CONFIG_DIR}/unbound.conf.example"
            bashio::log.info "Example config created at: <HA config>/unbound/unbound.conf.example"
            bashio::log.info "Rename it to unbound.conf and restart the add-on to use it"
        fi
    fi
else
    bashio::log.warning "/config is not available, using /data for config storage"
    CONFIG_DIR="/data"

    if [ ! -f "/data/unbound.conf" ] && [ ! -f "/data/unbound.conf.example" ]; then
        if [ -f "/etc/unbound/unbound.conf.example" ]; then
            cp /etc/unbound/unbound.conf.example /data/unbound.conf.example
            bashio::log.info "Example config created at: /data/unbound.conf.example"
        fi
    fi
fi

# Check if custom config exists (check multiple locations)
if [ -f "/config/unbound/unbound.conf" ]; then
    bashio::log.info "Using custom unbound.conf from /config/unbound folder"
    cp /config/unbound/unbound.conf /etc/unbound/unbound.conf
elif [ -f "/data/unbound.conf" ]; then
    bashio::log.info "Using custom unbound.conf from /data folder"
    cp /data/unbound.conf /etc/unbound/unbound.conf
else
    # Default simple config
    bashio::log.info "Using default configuration"
    cat > /etc/unbound/unbound.conf << 'EOF'
EOF
fi

# Validate configuration
bashio::log.info "Validating Unbound configuration..."
if ! unbound-checkconf /etc/unbound/unbound.conf; then
    bashio::log.error "Configuration validation failed!"
    bashio::log.error "Please check your unbound.conf file for errors"
    exit 1
fi

bashio::log.info "Configuration is valid, starting Unbound..."

# Start unbound
exec unbound -d -c /etc/unbound/unbound.conf
